# templates/backup.mount.j2
[Unit]
Description=Mount /backup

[Mount]
What=/dev/disk/by-uuid/{{ uuid_backup }}
Where=/backup
Type=btrfs
Options=defaults,rw,noatime,space_cache

[Install]
WantedBy=multi-user.target